# systemd/hyprcaffeine.service
[Unit]
Description=HyprCaffeine — restore inhibitors on user-session start

[Service]
Type=oneshot
# restore waits internally for the Hyprland socket (up to 60 s) before
# creating inhibitors, so set a generous timeout here.
TimeoutStartSec=120
ExecStart=%h/.local/bin/hyprcaffeine restore
ExecStart=-%h/.local/bin/hyprcaffeine watcher start
RemainAfterExit=yes

[Install]
WantedBy=default.target
