# systemd/hyprcaffeine.service
[Unit]
Description=HyprCaffeine — restore inhibitors on user-session start

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprcaffeine restore
ExecStart=%h/.local/bin/hyprcaffeine watcher start
RemainAfterExit=yes

[Install]
WantedBy=default.target
